# backend/deploy/energy-audit.service
[Unit]
Description=Energy Audit web app
After=network.target

[Service]
WorkingDirectory=/opt/energy-audit
EnvironmentFile=/opt/energy-audit/.env

# Bind only to localhost; nginx proxies to it.
ExecStart=/opt/energy-audit/.venv/bin/gunicorn -b 127.0.0.1:5050 wsgi:app

Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target

